# crontab
*/15 * * * * cd /nothilfe_tracker/src && /usr/local/bin/python tracker.py >> /var/log/cron.log 2>&1
